# govai-api.service
[Unit]
Description=GovAI API Service
After=syslog.target network.target

[Service]
User=ubuntu
Group=ubuntu

# Working directory
WorkingDirectory=/opt/govai/backend

# Environment file
EnvironmentFile=/opt/govai/backend/.env

# Python executable path
ExecStart=/opt/govai/venv/bin/gunicorn app.main:app \
    --workers 2 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 0.0.0.0:8000 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile -

# Restart policy
Restart=always
RestartSec=10
SuccessExitStatus=143

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=govai-api

# Security
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
