# govai-frontend.service
[Unit]
Description=GovAI Frontend Service
After=syslog.target network.target

[Service]
User=ubuntu
Group=ubuntu

# Working directory
WorkingDirectory=/opt/govai/frontend

# Environment variables
Environment=NODE_ENV=production
Environment=PORT=3005
Environment=HOSTNAME=0.0.0.0
Environment=NEXT_PUBLIC_API_URL=http://35.173.103.83:8000

# Node executable path
ExecStart=/usr/bin/node server.js

# Restart policy
Restart=always
RestartSec=10
SuccessExitStatus=143

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=govai-frontend

# Security
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
